# os/services/holmium-scheduler
#!/sbin/openrc-run
# holmium-scheduler — Background scheduler loop
# Runs scheduled tasks from /var/holmium/scheduler.json

description="Holmium OS — Background Scheduler"
pidfile="/run/holmium/scheduler.pid"
command="/usr/bin/python"
command_args="-c \"from backend.scheduler_runner import run; run()\""
command_user="holmium:holmium"
command_background=true

depend() {
    need net
    after holmium-backend
    use logger
}

start() {
    ebegin "Starting holmium-scheduler"
    checkpath -d -m 0755 -o holmium:holmium /run/holmium
    start-stop-daemon --start \
        --pidfile "$pidfile" \
        --user "$command_user" \
        --make-pidfile \
        --background \
        --exec "$command" -- $command_args
    eend $?
}

stop() {
    ebegin "Stopping holmium-scheduler"
    start-stop-daemon --stop --pidfile "$pidfile" --retry 10
    rm -f "$pidfile"
    eend $?
}
